# supervisord.conf
[program:node]
command = node server.js
environment = NODE_ENV=production
directory = /home/dotcloud/current
